# giveaway_bot.service
[Unit]
Description=Discord Giveaway Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/giveaway_bot
ExecStart=/opt/giveaway_bot/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sicherheit
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
